# resources/debian/postrm
#!/bin/sh
# postrm script for rndpwd

set -e

update-rc.d -f rndpwd remove

#DEBHELPER#

exit 0
